# systemd/typing_ex_normal.timer
[Unit]
Description=typing_ex stat=normal

[Timer]
#OnCalendar=*:*:00
OnCalendar=Wed 9:00
OnCalendar=Wed 10:40

[Install]
WantedBy=timers.target
