# systemd/typing_ex_normal.timer
[Unit]
Description=typing_ex stat=normal

[Timer]
OnCalendar=Wed *-*-* 9:00,10:40
OnCalendar=*-*-* *:*:30

[Install]
WantedBy=timers.target
