# systemd/typing_ex_roll_call.timer
[Unit]
Description=typing_ex stat=normal

[Timer]
OnCalendar=Wed *-*-* 8:45,10:25
OnCalendar=*-*-* *:*:00

[Install]
WantedBy=timers.target
